# scripts/systemd/rocket-blog-backup.timer
[Unit]
Description=Rocket Blog Docker Volume Backup Timer
Documentation=https://github.com/xiro-codes/rocket_blog
Requires=rocket-blog-backup.service

[Timer]
# Run daily at 2:00 AM
OnCalendar=daily
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target