# scripts/systemd/rocket-blog-backup.service
[Unit]
Description=Rocket Blog Docker Volume Backup
Documentation=https://github.com/xiro-codes/rocket_blog
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/rocket_blog
ExecStart=/opt/rocket_blog/scripts/docker-deploy.sh backup prod
ExecStartPost=/opt/rocket_blog/scripts/docker-deploy.sh backup-clean 7
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/rocket_blog/backups
PrivateTmp=true